# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/app/logs/supervisord.pid

[program:main_system]
command=python /app/src/main.py
autostart=true
autorestart=true
stderr_logfile=/app/logs/main.err.log
stdout_logfile=/app/logs/main.out.log
environment=PYTHONUNBUFFERED="1"
priority=1

[program:dashboard]
command=python /app/src/dashboard.py
autostart=true
autorestart=true
stderr_logfile=/app/logs/dashboard.err.log
stdout_logfile=/app/logs/dashboard.out.log
environment=PYTHONUNBUFFERED="1"
priority=2

[program:higgsfield_agent]
command=python /app/src/higgsfield_agent.py
autostart=true
autorestart=true
stderr_logfile=/app/logs/higgsfield.err.log
stdout_logfile=/app/logs/higgsfield.out.log
environment=PYTHONUNBUFFERED="1"
priority=3

[group:kimmy]
programs=main_system,dashboard,higgsfield_agent